# infra/systemd/normaai-photo-watcher.timer
[Unit]
Description=NormaAI · Trigger photo watcher every 5 minutes
Requires=normaai-photo-watcher.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=normaai-photo-watcher.service

[Install]
WantedBy=timers.target
